# .config/systemd/user/cycle-mac.service
[Unit]
Description=Check internet connectivity and cycle MAC if needed

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/cycle_mac
