# .config/systemd/user/cycle-mac.timer
[Unit]
Description=Run MAC cycle check every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
